# files/ghidra.service.in
[Unit]
Description=Ghidra Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=ghidra
Group=ghidra

# Pull out a bunch of stuff the ghidra server obviously won't need
# TODO: We should even restrict syscalls using systemd seccomp filtering
ProtectHome=yes
ProtectSystem=strict
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
LockPersonality=yes
DevicePolicy=closed
#MountAPIVFS=yes
PrivateDevices=yes
PrivateTmp=yes

ReadWritePaths=@ghidra_home@ @ghidra_root@

Environment=JAVA=@ghidra_java@
Environment=GHIDRA_HOME=@ghidra_home@
Environment=GHIDRA_ROOT=@ghidra_root@
Environment=GHIDRA_ARGS="-a0 -i @public_ip@ @anonymous@"

# this is a certified mess. Thanks lennart :)
ExecStart=@ghidra_root@/ghidra_systemd_run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
