# files/ghidra.service.in
[Unit]
Description=Ghidra Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=ghidra
Group=ghidra

# Misc service hardening options

RestrictRealtime=yes
# Ghidra really only needs INET, ipv6 kinda works but it's a bit janky
# and ghidra tends to freak out if you use it
RestrictAddressFamilies=AF_INET
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=no

NoNewPrivileges=yes

MountAPIVFS=yes
DevicePolicy=closed
PrivateDevices=yes
PrivateTmp=yes

BindReadOnlyPaths=@ghidra_home@
BindPaths=@ghidra_root@

ProtectHome=yes
ProtectSystem=strict
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectProc=invisible
LockPersonality=yes

Environment=JAVA=@ghidra_java@
Environment=GHIDRA_HOME=@ghidra_home@
Environment=GHIDRA_ROOT=@ghidra_root@
Environment=GHIDRA_ARGS="-a0 -i @public_ip@ @anonymous@"

ExecStart=@ghidra_root@/ghidra_systemd_run
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
